# templates/etc/supervisor/service.conf.erb
[program:<%= @title %>]
command=<%= @command %>
directory=<%= @directory %>
user=<%= @user %>
process_name=<%= @process_name %>
startsecs=<%= @startsecs %>
numprocs=<%= @numprocs %>
autostart=<%= @autostart %>
autorestart=<%= @autorestart %>
redirect_stderr=<%= @redirect_stderr %>
startretries=<%= @startretries %>
